# systemd/ansible-estate.timer
[Unit]
Description=Push the ZMS estate with orchestrate.yml (~every 60 minutes)

[Timer]
OnBootSec=15min
OnUnitActiveSec=60min
RandomizedDelaySec=120
Unit=ansible-estate.service

[Install]
WantedBy=timers.target
