# systemd/ansible-estate.service
[Unit]
Description=ZMS Ansible estate push (orchestrate.yml - full estate)
After=network-online.target ansible-bootstrap.service
Wants=network-online.target

[Service]
Type=oneshot
User=ansible
Group=ansible
EnvironmentFile=/etc/ansible/estate.env
WorkingDirectory=/opt/control-repo
ExecStart=/opt/ansible-venv/bin/ansible-playbook orchestrate.yml
ExecStopPost=/bin/sh -c 'bash /opt/control-repo/scripts/notify-result.sh "$EXIT_STATUS" ansible-estate'
TimeoutStartSec=7200
